# buildroot/board/rockchip/px3se/fs-overlay/etc/profile.d/RkUpdater.sh
#! /bin/bash
sleep 3
echo -e "\n----->Start updater <-----\n"
updater -d